# build/Lede_source/files/etc/crontabs/root
* * * * * sh /etc/config/my_watchdog.sh
#0 */1 * * * /etc/init.d/ddnsto restart
0 1 * * 1 rm /tmp/luci-indexcache > /dev/null 2>&1
