# build/Lede_source/files/etc/crontabs/root
*/40 * * * * /etc/init.d/frp restart
* * * * * sh /etc/config/my_watchdog.sh
