# build/Lede_source/files/etc/crontabs/root
* * * * * sh /etc/config/my_watchdog.sh
0 */1 * * * /etc/init.d/ddnsto restart
*/40 * * * * /etc/init.d/frp restart
